# gentoo/overlay/sys-apps/systemd/files/systemd-units/vixie-cron.service
[Unit]
Description=Vixie Cron Daemon
Wants=syslog.target

[Service]
ExecStart=/usr/sbin/cron -n
Restart=always

[Install]
WantedBy=multi-user.target
